# systemd/custom-init.service
[Unit]
Description=CustomInit

[Service]
ExecStart=/tmp/custom-init.sh

[Install]
WantedBy=multi-user.target